# bios-configd.service
[Unit]
Description=BIOS settings manager

[Service]
Restart=always
ExecStart=/usr/bin/bios-configd

[Install]
WantedBy=multi-user.target
